# debian/cocaine-isolate-daemon-upstart.upstart
description     "cocaine isolate daemon"

start on filesystem
stop on runlevel [06]

kill timeout 30
respawn

script
    exec cocaine-isolate-daemon --config /etc/cocaine-isolate-daemon/cocaine-isolate-daemon.conf
end script
